# analytics/CAR-2020-08-003.yaml
---
title: Scheduled Task - FileAccess
submission_date: 2020/08/31
information_domain: 'Host'
platforms:
  - Windows
subtypes: 
  - File
analytic_types:
  - TTP
contributors:
- Olaf Hartong
id: CAR-2020-08-003
description: |-
  Description text goes here.
coverage:
- coverage: Low
  subtechniques: []
  tactics: []
  technique: T1053
implementations:
- name: 
  description: ''
  code: |- 
    index=__your_sysmon_index__ event_id=11 process_path!="C:\\\\WINDOWS\\\\\system32\\\\svchost.exe" (file_path="C:\\\\Windows\\\\System32\\\\Tasks\\\\\
    *\" OR file_path=\"C:\\\\Windows\\\\Tasks\\\\*")
  data_model: Sysmon native
  type: Splunk

