# analytics/CAR-2013-01-002.yaml
title: Autorun Differences
submission_date: 2013/01/25
information_domain: 'Analytic, Host'
subtypes:
  - Registry
analytic_types:
  - Situational Awareness
  - TTP
contributors:
  - MITRE
id: CAR-2013-01-002
description: |-
  The Sysinternals tool [Autoruns](../sensors/autoruns) checks the registry and file system for known identify persistence mechanisms. It will output any tools identified, including built-in or added-on Microsoft functionality and third party software. Many of these locations are known by adversaries and used to obtain [Persistence](https://attack.mitre.org/tactics/TA0003). Running Autoruns periodically in an environment makes it possible to collect and monitor its output for differences, which may include the removal or addition of persistent tools. Depending on the persistence mechanism and location, legitimate software may be more likely to make changes than an adversary tool. Thus, this analytic may result in significant noise in a highly dynamic environment. While Autoruns is a convenient method to scan for programs using persistence mechanisms its scanning nature does not conform well to streaming based analytics. This analytic could be replaced with one that draws from sensors that collect registry and file information if streaming analytics are desired.

  Utilizes the Sysinternals autoruns tool (ignoring validated Microsoft entries). Primarily not a detection analytic by itself but through analysis of results by an analyst can be used for such. Building another analytic on top of this one identifying unusual entries would likely be a beneficial alternative.
coverage:
  - technique: T1031
    tactics:
      - TA0003
    coverage: Moderate
  - technique: T1050
    tactics:
      - TA0003
    coverage: Moderate
  - technique: T1053
    tactics:
      - TA0003
    coverage: Moderate
  - technique: T1013
    tactics:
      - TA0003
    coverage: Moderate
  - technique: T1060
    tactics:
      - TA0003
    coverage: Moderate
  - technique: T1034
    tactics:
      - TA0003
      - TA0004
    coverage: Moderate
  - technique: T1015
    tactics:
      - TA0004
      - TA0003
    coverage: Moderate
  - technique: T1112
    tactics:
      - TA0003
      - TA0002
    coverage: Moderate
  - technique: T1058
    tactics:
      - TA0003
      - TA0002
    coverage: Moderate
  - technique: T1084
    tactics:
      - TA0003
    coverage: Moderate
  - technique: T1044
    tactics:
      - TA0004
      - TA0003
    coverage: Moderate
  - technique: T1042
    tactics:
      - TA0003
    coverage: Moderate
  - technique: T1037
    tactics:
      - TA0003
    coverage: Moderate
  - technique: T1004
    tactics:
      - TA0003
    coverage: Moderate
  - technique: T1103
    tactics:
      - TA0003
    coverage: Moderate
