# docs/car_attack/car_attack.json
{
    "domain": "mitre-enterprise", 
    "version": "2.1", 
    "techniques": [
        {
            "color": "#c6dbef", 
            "comment": "CAR-2013-01-003: SMB Events Monitoring | CAR-2013-02-008: Simultaneous Logins on a Host | CAR-2013-02-012: User Logged in to Multiple Hosts | CAR-2013-05-003: SMB Write Request | CAR-2013-05-005: SMB Copy and Execution | CAR-2013-10-001: User Login Activity Monitoring | CAR-2014-05-001: RPC Activity | CAR-2016-04-005: Remote Desktop Logon", 
            "techniqueID": "T1078", 
            "enabled": true
        }, 
        {
            "color": "#c6dbef", 
            "comment": "CAR-2013-01-002: Autorun Differences | CAR-2014-11-003: Debuggers for Accessibility Applications | CAR-2014-11-003: Debuggers for Accessibility Applications | CAR-2014-11-008: Command Launched from WinLogon", 
            "techniqueID": "T1015", 
            "enabled": true
        }, 
        {
            "color": "#c6dbef", 
            "comment": "CAR-2013-01-002: Autorun Differences", 
            "techniqueID": "T1037", 
            "enabled": true
        }, 
        {
            "color": "#c6dbef", 
            "comment": "CAR-2013-05-002: Suspicious Run Locations | CAR-2013-05-009: Running executables with same hash and different names | CAR-2013-07-001: Suspicious Arguments", 
            "techniqueID": "T1036", 
            "enabled": true
        }, 
        {
            "color": "#c6dbef", 
            "comment": "CAR-2013-04-002: Quick execution of a series of suspicious commands | CAR-2014-02-001: Service Binary Modifications | CAR-2014-03-005: Remotely Launched Executables via Services", 
            "techniqueID": "T1035", 
            "enabled": true
        }, 
        {
            "color": "#c6dbef", 
            "comment": "CAR-2013-01-002: Autorun Differences | CAR-2014-07-001: Service Search Path Interception", 
            "techniqueID": "T1034", 
            "enabled": true
        }, 
        {
            "color": "#c6dbef", 
            "comment": "CAR-2013-01-002: Autorun Differences | CAR-2013-04-002: Quick execution of a series of suspicious commands | CAR-2013-09-005: Service Outlier Executables | CAR-2014-02-001: Service Binary Modifications | CAR-2014-03-005: Remotely Launched Executables via Services", 
            "techniqueID": "T1031", 
            "enabled": true
        }, 
        {
            "color": "#c6dbef", 
            "comment": "CAR-2019-04-002: Generic Regsvr32 | CAR-2019-04-003: Squiblydoo", 
            "techniqueID": "T1117", 
            "enabled": true
        }, 
        {
            "color": "#c6dbef", 
            "comment": "CAR-2013-04-002: Quick execution of a series of suspicious commands", 
            "techniqueID": "T1063", 
            "enabled": true
        }, 
        {
            "color": "#c6dbef", 
            "comment": "CAR-2013-01-003: SMB Events Monitoring", 
            "techniqueID": "T1039", 
            "enabled": true
        }, 
        {
            "color": "#c6dbef", 
            "comment": "CAR-2013-10-002: DLL Injection via Load Library | CAR-2013-10-002: DLL Injection via Load Library", 
            "techniqueID": "T1055", 
            "enabled": true
        }, 
        {
            "color": "#c6dbef", 
            "comment": "CAR-2013-04-002: Quick execution of a series of suspicious commands", 
            "techniqueID": "T1018", 
            "enabled": true
        }, 
        {
            "color": "#c6dbef", 
            "comment": "CAR-2013-04-002: Quick execution of a series of suspicious commands | CAR-2016-03-001: Host Discovery Commands", 
            "techniqueID": "T1057", 
            "enabled": true
        }, 
        {
            "color": "#c6dbef", 
            "comment": "CAR-2013-01-003: SMB Events Monitoring | CAR-2013-04-002: Quick execution of a series of suspicious commands | CAR-2013-05-003: SMB Write Request | CAR-2013-05-005: SMB Copy and Execution", 
            "techniqueID": "T1077", 
            "enabled": true
        }, 
        {
            "color": "#c6dbef", 
            "comment": "CAR-2013-01-002: Autorun Differences | CAR-2013-04-002: Quick execution of a series of suspicious commands | CAR-2013-09-005: Service Outlier Executables | CAR-2014-02-001: Service Binary Modifications | CAR-2014-03-005: Remotely Launched Executables via Services | CAR-2014-05-001: RPC Activity | CAR-2014-05-002: Services launching Cmd", 
            "techniqueID": "T1050", 
            "enabled": true
        }, 
        {
            "color": "#c6dbef", 
            "comment": "CAR-2013-01-002: Autorun Differences | CAR-2013-04-002: Quick execution of a series of suspicious commands | CAR-2013-05-004: Execution with AT | CAR-2013-08-001: Execution with schtasks | CAR-2014-05-001: RPC Activity | CAR-2015-04-001: Remotely Scheduled Tasks via AT | CAR-2015-04-002: Remotely Scheduled Tasks via Schtasks", 
            "techniqueID": "T1053", 
            "enabled": true
        }, 
        {
            "color": "#c6dbef", 
            "comment": "CAR-2016-04-004: Successful Local Account Login | CAR-2016-04-004: Successful Local Account Login", 
            "techniqueID": "T1075", 
            "enabled": true
        }, 
        {
            "color": "#c6dbef", 
            "comment": "CAR-2013-04-002: Quick execution of a series of suspicious commands", 
            "techniqueID": "T1010", 
            "enabled": true
        }, 
        {
            "color": "#c6dbef", 
            "comment": "CAR-2013-01-002: Autorun Differences", 
            "techniqueID": "T1013", 
            "enabled": true
        }, 
        {
            "color": "#c6dbef", 
            "comment": "CAR-2013-03-001: Reg.exe called from Command Shell | CAR-2013-04-002: Quick execution of a series of suspicious commands", 
            "techniqueID": "T1012", 
            "enabled": true
        }, 
        {
            "color": "#c6dbef", 
            "comment": "CAR-2013-02-003: Processes Spawning cmd.exe | CAR-2013-04-002: Quick execution of a series of suspicious commands | CAR-2014-11-002: Outlier Parents of Cmd", 
            "techniqueID": "T1059", 
            "enabled": true
        }, 
        {
            "color": "#c6dbef", 
            "comment": "CAR-2013-01-002: Autorun Differences | CAR-2013-03-001: Reg.exe called from Command Shell | CAR-2013-04-002: Quick execution of a series of suspicious commands", 
            "techniqueID": "T1058", 
            "enabled": true
        }, 
        {
            "color": "#c6dbef", 
            "comment": "CAR-2013-01-002: Autorun Differences | CAR-2013-03-001: Reg.exe called from Command Shell | CAR-2013-04-002: Quick execution of a series of suspicious commands | CAR-2014-11-005: Remote Registry", 
            "techniqueID": "T1112", 
            "enabled": true
        }, 
        {
            "color": "#c6dbef", 
            "comment": "CAR-2013-04-002: Quick execution of a series of suspicious commands | CAR-2016-03-001: Host Discovery Commands", 
            "techniqueID": "T1016", 
            "enabled": true
        }, 
        {
            "color": "#c6dbef", 
            "comment": "CAR-2013-04-002: Quick execution of a series of suspicious commands", 
            "techniqueID": "T1098", 
            "enabled": true
        }, 
        {
            "color": "#c6dbef", 
            "comment": "CAR-2013-04-002: Quick execution of a series of suspicious commands | CAR-2016-03-001: Host Discovery Commands", 
            "techniqueID": "T1082", 
            "enabled": true
        }, 
        {
            "color": "#c6dbef", 
            "comment": "CAR-2013-04-002: Quick execution of a series of suspicious commands", 
            "techniqueID": "T1046", 
            "enabled": true
        }, 
        {
            "color": "#c6dbef", 
            "comment": "CAR-2014-04-003: Powershell Execution | CAR-2014-04-003: Powershell Execution | CAR-2014-04-003: Powershell Execution | CAR-2014-11-004: Remote PowerShell Sessions | CAR-2014-11-004: Remote PowerShell Sessions | CAR-2015-04-002: Remotely Scheduled Tasks via Schtasks", 
            "techniqueID": "T1086", 
            "enabled": true
        }, 
        {
            "color": "#c6dbef", 
            "comment": "CAR-2013-04-002: Quick execution of a series of suspicious commands | CAR-2016-03-001: Host Discovery Commands", 
            "techniqueID": "T1087", 
            "enabled": true
        }, 
        {
            "color": "#c6dbef", 
            "comment": "CAR-2013-01-002: Autorun Differences", 
            "techniqueID": "T1084", 
            "enabled": true
        }, 
        {
            "color": "#c6dbef", 
            "comment": "CAR-2013-04-002: Quick execution of a series of suspicious commands | CAR-2016-04-002: User Activity from Clearing Event Logs | CAR-2016-04-003: User Activity from Stopping Windows Defensive Services", 
            "techniqueID": "T1054", 
            "enabled": true
        }, 
        {
            "color": "#c6dbef", 
            "comment": "CAR-2013-10-002: DLL Injection via Load Library | CAR-2019-04-001: UAC Bypass", 
            "techniqueID": "T1088", 
            "enabled": true
        }, 
        {
            "color": "#c6dbef", 
            "comment": "CAR-2013-04-002: Quick execution of a series of suspicious commands", 
            "techniqueID": "T1089", 
            "enabled": true
        }, 
        {
            "color": "#c6dbef", 
            "comment": "CAR-2013-07-001: Suspicious Arguments | CAR-2014-05-001: RPC Activity", 
            "techniqueID": "T1021", 
            "enabled": true
        }, 
        {
            "color": "#c6dbef", 
            "comment": "CAR-2014-03-006: RunDLL32.exe monitoring | CAR-2014-03-006: RunDLL32.exe monitoring", 
            "techniqueID": "T1085", 
            "enabled": true
        }, 
        {
            "color": "#c6dbef", 
            "comment": "CAR-2014-11-004: Remote PowerShell Sessions | CAR-2014-11-006: Windows Remote Management (WinRM) | CAR-2014-11-006: Windows Remote Management (WinRM)", 
            "techniqueID": "T1028", 
            "enabled": true
        }, 
        {
            "color": "#c6dbef", 
            "comment": "CAR-2013-04-002: Quick execution of a series of suspicious commands", 
            "techniqueID": "T1029", 
            "enabled": true
        }, 
        {
            "color": "#c6dbef", 
            "comment": "CAR-2013-07-002: RDP Connection Detection | CAR-2013-07-002: RDP Connection Detection | CAR-2013-10-001: User Login Activity Monitoring | CAR-2014-11-008: Command Launched from WinLogon | CAR-2016-04-005: Remote Desktop Logon", 
            "techniqueID": "T1076", 
            "enabled": true
        }, 
        {
            "color": "#c6dbef", 
            "comment": "CAR-2013-01-002: Autorun Differences", 
            "techniqueID": "T1103", 
            "enabled": true
        }, 
        {
            "color": "#c6dbef", 
            "comment": "CAR-2014-11-007: Remote Windows Management Instrumentation (WMI) over RPC | CAR-2014-11-007: Remote Windows Management Instrumentation (WMI) over RPC | CAR-2014-12-001: Remotely Launched Executables via WMI | CAR-2016-03-002: Create Remote Process via WMIC | CAR-2016-03-002: Create Remote Process via WMIC", 
            "techniqueID": "T1047", 
            "enabled": true
        }, 
        {
            "color": "#c6dbef", 
            "comment": "CAR-2013-01-002: Autorun Differences | CAR-2014-02-001: Service Binary Modifications", 
            "techniqueID": "T1044", 
            "enabled": true
        }, 
        {
            "color": "#c6dbef", 
            "comment": "CAR-2013-01-002: Autorun Differences", 
            "techniqueID": "T1042", 
            "enabled": true
        }, 
        {
            "color": "#c6dbef", 
            "comment": "CAR-2013-04-002: Quick execution of a series of suspicious commands | CAR-2016-03-001: Host Discovery Commands", 
            "techniqueID": "T1069", 
            "enabled": true
        }, 
        {
            "color": "#c6dbef", 
            "comment": "CAR-2013-05-003: SMB Write Request | CAR-2013-05-005: SMB Copy and Execution | CAR-2013-07-001: Suspicious Arguments", 
            "techniqueID": "T1105", 
            "enabled": true
        }, 
        {
            "color": "#c6dbef", 
            "comment": "CAR-2013-07-005: Command Line Usage of Archiving Software", 
            "techniqueID": "T1002", 
            "enabled": true
        }, 
        {
            "color": "#c6dbef", 
            "comment": "CAR-2013-04-002: Quick execution of a series of suspicious commands | CAR-2013-07-001: Suspicious Arguments | CAR-2019-04-004: Credential Dumping via Mimikatz", 
            "techniqueID": "T1003", 
            "enabled": true
        }, 
        {
            "color": "#c6dbef", 
            "comment": "CAR-2013-04-002: Quick execution of a series of suspicious commands | CAR-2016-03-001: Host Discovery Commands", 
            "techniqueID": "T1007", 
            "enabled": true
        }, 
        {
            "color": "#c6dbef", 
            "comment": "CAR-2013-01-002: Autorun Differences", 
            "techniqueID": "T1004", 
            "enabled": true
        }, 
        {
            "color": "#c6dbef", 
            "comment": "CAR-2013-04-002: Quick execution of a series of suspicious commands", 
            "techniqueID": "T1049", 
            "enabled": true
        }, 
        {
            "color": "#c6dbef", 
            "comment": "CAR-2014-04-003: Powershell Execution", 
            "techniqueID": "T1064", 
            "enabled": true
        }
    ], 
    "name": "CAR ATT&CK", 
    "description": "CAR Analytics ATT&CK Coverage"
}